# crontab
SHELL=/bin/bash
20 0 * * * /usr/local/bin/curator delete indices --older-than 7 --time-unit days --timestring "\%Y.\%m.\%d"
